# addon-ha-spawn-point-bedrock/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
DOWNLOAD_SAVE_DIR=/data/download
INSTALL_DIR=/data/bedrock_server
RUN_VERSION=$(bashio::config 'version')

# Stores the location of this library
readonly __HASP_LIB_DIR=$(dirname "${BASH_SOURCE[0]}")
source "${__HASP_LIB_DIR}/hasp_config.sh"

## Print banner
./banner.sh

## EULA check
EULA=$(bashio::config 'eula')
if [[ ${EULA^^} != TRUE ]]; then
  bashio::log.red "EULA must be set to TRUE to indicate agreement with the Minecraft End User License"
  bashio::log.red "See https://minecraft.net/terms"
  bashio::log.red "Current value is '${EULA}'"
  exit 1
fi

if [[ "$RUN_VERSION" = "LATEST" ]]; then
  curl https://mc-bds-helper.vercel.app/api/latest -o /data/latest_version_download_link.txt > /dev/null 2>&1
  if [[ -f "/data/latest_version_download_link.txt" ]]; then
    DOWNLOAD_URL=$(cat /data/latest_version_download_link.txt)
    bashio::log.green "Download URL '${DOWNLOAD_URL}'"
  fi
fi

bashio::log.green "Download Directory ${DOWNLOAD_SAVE_DIR}"

if [[ ! -d  ${DOWNLOAD_SAVE_DIR} ]]; then
  bashio::log "Creating Download Directory"
  mkdir ${DOWNLOAD_SAVE_DIR}
  chmod -R 777 ${DOWNLOAD_SAVE_DIR}
fi

DOWNLOAD_FILE=bedrock-server-1.20.41.02.zip

if [[ ! -f "${DOWNLOAD_SAVE_DIR}/${DOWNLOAD_FILE}" ]]; then
  bashio::log.yellow "Download Verison %{RUN_VERSION} started ... please wait"
  curl "https://minecraft.azureedge.net/bin-linux/${DOWNLOAD_FILE}" -o "${DOWNLOAD_SAVE_DIR}/${DOWNLOAD_FILE}" -# > /dev/null 2>&1
else
  bashio::log.green "${DOWNLOAD_FILE} already downloaded"
fi

if [[ ! -d  ${INSTALL_DIR} ]]; then
  bashio::log "Creating Install Directory"
  mkdir ${INSTALL_DIR}
  chmod -R 777 ${INSTALL_DIR}
fi

SERVER_FILE_NAME=bedrock_server-1.20.41.02

if [[ ! -f "${INSTALL_DIR}/${SERVER_FILE_NAME}" ]]; then
  bashio::log "Installing new version"
  unzip -q -n "${DOWNLOAD_SAVE_DIR}/${DOWNLOAD_FILE}" -d ${INSTALL_DIR}
  mv "${INSTALL_DIR}/bedrock_server" "${INSTALL_DIR}/${SERVER_FILE_NAME}"
else
  bashio::log.green "${DOWNLOAD_FILE} already Installed"
fi

cd "${INSTALL_DIR}" || exit 1

#create server configuration based on ha ui configuration
bashio::log "Create Bedrock Server Configuration"
hasp::config.update_server server-name server-name ./server.properties
hasp::config.update_server level_name level-name ./server.properties
hasp::config.update_server level_gamemode gamemode ./server.properties
hasp::config.update_server level_difficulty difficulty ./server.properties
hasp::config.update_server level_seed level-seed ./server.properties




#prepare and start bedrock server
bashio::log "Starting Bedrock Server"
chmod a+x "./${SERVER_FILE_NAME}"
export LD_LIBRARY_PATH=.
"./${SERVER_FILE_NAME}"



#exec python3 /usr/bin/mc_port_proxy.py
## Run your program
#exec /usr/bin/my_program
