# addon-ha-spawn-point-bedrock/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare message

DOWNLOAD_SAVE_DIR=/data/download
INSTALL_DIR=/data/bedrock_server
#DOWNLOAD_SAVE_DIR+=addon_spawn_point

## EULA check
EULA=$(bashio::config 'eula')

if [[ ${EULA^^} != TRUE ]]; then
  bashio::log.red "EULA must be set to TRUE to indicate agreement with the Minecraft End User License"
  bashio::log.red "See https://minecraft.net/terms"
  bashio::log.red "Current value is '${EULA}'"
  exit 1
fi

## Print banner
./banner.sh


## Get the 'server_name' key from the user config options.
#message=$(bashio::config 'server_name')

## Print the message the user supplied, defaults to "Hello World..."
#bashio::log.info "${message:="Hello World..."}"
bashio::log.green "Download Directory ${DOWNLOAD_SAVE_DIR}"

if [[ ! -d  ${DOWNLOAD_SAVE_DIR} ]]; then
  bashio::log "Creating Download Directory"
  mkdir ${DOWNLOAD_SAVE_DIR}
  chmod -R 666 ${DOWNLOAD_SAVE_DIR}
fi

DOWNLOAD_FILE=bedrock-server-1.20.41.02.zip

if [[ ! -f "${DOWNLOAD_SAVE_DIR}/${DOWNLOAD_FILE}" ]]; then
  wget "https://minecraft.azureedge.net/bin-linux/${DOWNLOAD_FILE}" -P ${DOWNLOAD_SAVE_DIR}/
else
  bashio::log.green "${DOWNLOAD_FILE} already downloaded"
fi

if [[ ! -d  ${INSTALL_DIR} ]]; then
  bashio::log "Creating Install Directory"
  mkdir ${INSTALL_DIR}
  chmod -R 666 ${DOWNLOAD_SAVE_DIR}
fi

if [[ ! -f "${INSTALL_DIR}/bedrock-server-1.20.41.02" ]]; then
  bashio::log "Installing new version"
  unzip -q -n "${DOWNLOAD_SAVE_DIR}/${DOWNLOAD_FILE}" -d ${INSTALL_DIR}
  mv "${INSTALL_DIR}/bedrock_server" "${INSTALL_DIR}/bedrock_server-1.20.41.02"
else
  bashio::log "Already Installed"
fi


exec python3 /usr/bin/mc_port_proxy.py
## Run your program
#exec /usr/bin/my_program
